# letsencrypt/bin/crontab
# m h d mth wd cmd
26 4 15 * * /opt/letsencrypt/bin/update-certs.sh
